# clockmonitor.service
[Unit]
Description=ClockMonitor FastAPI dashboard
After=network-online.target chrony.service gpsd.service ptp4l.service
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/ClockMonitor
ExecStart=/home/YOUR_USER/ClockMonitor/.venv/bin/uvicorn app.main:create_app --factory --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3
PrivateTmp=false
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
